# exgit@.service
[Unit]
Description=A job that enables exgit sync its repositories


[Service]
Environment="MOUNTPATH=%I"
