# exgit_sync.conf
description "A job that enables exgit sync its repositories "
author "Oginga Steven"

start on runlevel [2345]

instance $P
console log